# roles/flask_app/templates/supervisor_flask_app.conf.j2
[program:flask_app]
command = bash run_app.bash run
directory = {{ app_dir }} 
user = vagrant
